# ansible/roles/deploy/templates/systemd/rpims-watcher.path.j2
[Unit]
Description="Monitor the /opt/RPiMS/config/rpims.yaml file for changes"

[Path]
PathModified=/opt/RPiMS/config/rpims.yaml
Unit=rpims-watcher.service

[Install]
WantedBy=multi-user.target
