# ansible/roles/deploy/templates/rpims/rpims-watcher.path.j2
[Unit]
Description="Monitor the {{ rpims_install_dir }}/config/rpims.yaml file for changes"

[Path]
PathModified={{ rpims_install_dir }}/config/rpims.yaml

[Install]
WantedBy=multi-user.target
