# ansible/roles/deploy/templates/zabbix/zabbix-watcher.service.j2
[Unit]
Description=Zabbix restarter
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart zabbix-agent.service

[Install]
WantedBy=multi-user.target
